# .config/systemd/user/org-backup.service
[Unit]
Description=Org Backup Script

[Service]
Type=oneshot
ExecStart=/home/cleo/scripts/org-backup.sh
